# debian/DEBIAN/postinst
#!/bin/sh

if [ -x "/etc/init.d/nginx" ]; then
    update-rc.d nginx defaults >/dev/null
    /etc/init.d/nginx restart
fi
